# systemd/dehydrated.timer
[Unit]
Description=Renew searx' let's encrypt certificate with dehydrated

[Timer]
OnCalendar=daily
Persistent=true     

[Install]
WantedBy=timers.target
