# systemd/searx.service
[Unit]
Description=Searx
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
ExecStartPre=/usr/local/bin/docker-compose -f /var/searx-docker/docker-compose.yml kill
ExecStartPre=/usr/local/bin/docker-compose -f /var/searx-docker/docker-compose.yml rm -f
ExecStartPre=/usr/local/bin/docker-compose -f /var/searx-docker/docker-compose.yml pull
ExecStart=/usr/local/bin/docker-compose -f /var/searx-docker/docker-compose.yml up --remove-orphans --force-recreate --build
ExecStop=/usr/local/bin/docker-compose -f /var/searx-docker/docker-compose.yml stop

[Install]
WantedBy=multi-user.target
